# examples/nginfier.service
[Unit]
Description=Nginfier Service
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/nginfier
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings (optional but recommended)
ProtectSystem=full
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target 